# packaging/postinstall.sh
#!/usr/bin/env bash
systemctl enable dockerhub-pull-limit-exporter
systemctl start dockerhub-pull-limit-exporter
