# packaging/preremove.sh
#!/usr/bin/env bash
systemctl stop dockerhub-pull-limit-exporter
systemctl disable dockerhub-pull-limit-exporter
